# templates/curl.service.j2
[Unit]
Description=curl

[Service]
ExecStart=-/bin/curl -m 2 -s http://{{ hostvars['snort']['private_ip2'] }}/ddos_simulation
