# files/curl.timer
[Unit]
Description=Run curl every 5 seconds

[Timer]
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
