# files/ddosattack.timer
[Unit]
Description=Run ddosattack every 5 seconds

[Timer]
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
