# files/sqlinjection.timer
[Unit]
Description=Run sqlinjection every 5 seconds

[Timer]
OnUnitActiveSec=5s
AccuracySec=1s

[Install]
WantedBy=timers.target
